# units/pushbullet-test.timer
[Unit]
Description=Run pushbullet test every few minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
